# actionnetwork.timer
[Unit]
Description=Run daily subscriber update

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
