# actionnetwork.timer
[Unit]
Description=Run hourly subscriber update

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
